# etc/s6/nginx/run
#!/bin/sh

exec s6-setuidgid puppet nginx -g "daemon off;" -c /etc/nginx.conf
