# etc/s6/puma/run
#!/bin/sh

exec s6-setuidgid puppet puma --config ../../puppet/puma.conf
